# scripts/cursor-agent-sync.cron.example
# Install as a system cron snippet (runs as root — no sudo password, full /opt access).
#
#   sudo cp scripts/cursor-agent-sync.cron.example /etc/cron.d/formatforge-cursor-agent-sync
#   sudo chmod 644 /etc/cron.d/formatforge-cursor-agent-sync
#
# IMPORTANT: In cron.d, each job line must be EXACTLY ONE LINE. If your editor wraps
# long lines, the env vars or ">>/var/log/..." will break (cron is not shell multiline).
# Prefer the short line below that calls the wrapper script.
#
#   sudo cp scripts/cursor-agent-sync-root-invoke.sh.example /usr/local/sbin/formatforge-cursor-agent-sync
#   sudo nano /usr/local/sbin/formatforge-cursor-agent-sync && sudo chmod 755 /usr/local/sbin/formatforge-cursor-agent-sync

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Recommended: one short line (logging is inside the wrapper — works for manual runs too)
# m h dom mon dow user command
0 4 * * * root /usr/local/sbin/formatforge-cursor-agent-sync

# Alternative (all env on one line — verify with `cat -A` that nothing wrapped):
# 0 4 * * * root env CURSOR_AGENT_SOURCE_USER=jhs CURSOR_AGENT_OPT_DIR=/opt/cursor-agent CURSOR_AGENT_SYNC_RELOAD_PHP=1 /var/www/formatforge/scripts/sync-cursor-agent-to-opt.sh >>/var/log/formatforge-cursor-agent-sync.log 2>&1
